# crontab
m   h   d   m   y   cmd

5   *   *   *   *   php artisan parse:start
